# systemd/rss-sync.service
[Unit]
Description=Miniflux RSS Sync Service
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/rss-sync
ExecStart=/opt/rss-sync/venv/bin/python /opt/rss-sync/sync_miniflux.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
